# apiserver/supervisord.conf
[supervisord]
nodaemon=true

[program:celery]
command=celery -A plane worker -l info
autostart=true
autorestart=true
stderr_logfile=/var/log/celery.err.log
stdout_logfile=/var/log/celery.out.log

[program:flower]
command=celery -A plane flower --port=5555
autostart=true
autorestart=true
stderr_logfile=/var/log/flower.err.log
stdout_logfile=/var/log/flower.out.log