# deployments/aio/commercial/supervisor.conf
[supervisord]
user=root
nodaemon=true
stderr_logfile=/app/logs/error/supervisor.err.log
stdout_logfile=/app/logs/access/supervisor.log

[program:monitor]
directory=/app/monitor
command=sh -c "prime-monitor start"
autostart=true
autorestart=true
stdout_logfile_maxbytes=0
stdout_logfile=/app/logs/access/monitor.log
stderr_logfile=/app/logs/error/monitor.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3007,HOSTNAME=0.0.0.0
priority=10

[program:migrator]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-migrator.sh"
autostart=true
autorestart=unexpected
stdout_logfile=/app/logs/access/migrator.log
stderr_logfile=/app/logs/error/migrator.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10


[program:web]
command=sh -c "node /app/web/apps/web/server.js"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/web.log
stderr_logfile=/app/logs/error/web.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3001,HOSTNAME=0.0.0.0
priority=15

[program:space]
command=sh -c "node /app/space/apps/space/server.js"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/space.log
stderr_logfile=/app/logs/error/space.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3002,HOSTNAME=0.0.0.0
priority=15

[program:admin]
command=sh -c "node /app/admin/apps/admin/server.js"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/admin.log
stderr_logfile=/app/logs/error/admin.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3003,HOSTNAME=0.0.0.0
priority=15

[program:api]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-api.sh"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/api.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/api.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3004,HOSTNAME=0.0.0.0
priority=15


[program:worker]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-worker.sh"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/worker.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/worker.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:beat]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-beat.sh"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/beat.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/beat.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20


[program:live]
directory=/app/live
command=sh -c "node live/start.js"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/live.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/live.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3005,HOSTNAME=0.0.0.0
priority=20

[program:silo]
directory=/app/silo
command=sh -c "node silo/start.cjs"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/silo.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/silo.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=3006,HOSTNAME=0.0.0.0
priority=20

[program:email]
directory=/app/email
command=sh -c "./email-server"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/email.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/email.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:automation-consumer]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-automation-consumer.sh"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/automation-consumer.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/automation-consumer.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:outbox-poller]
directory=/app/backend
command=sh -c "./bin/docker-entrypoint-outbox-poller.sh"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/outbox-poller.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/outbox-poller.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:iframely]
directory=/app/iframely
command=sh -c "node server.js"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/iframely.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/iframely.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT=8061,HOSTNAME=0.0.0.0
priority=20

[program:proxy]
directory=/app/proxy
command=sh -c "caddy run --config /app/proxy/Caddyfile"
autostart=true
autorestart=true
stdout_logfile=/app/logs/access/proxy.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/logs/error/proxy.err.log
# stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20