# deploy/packer/svc/admin.service
[Unit]
Description=Plane Admin Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/admin
EnvironmentFile=/opt/plane/admin/admin.env
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && nvm use 20 && node admin/server.js'
Restart=always
RestartSec=10

# Logging configuration
StandardOutput=append:/var/log/plane/admin.log
StandardError=append:/var/log/plane/admin.error.log
SyslogIdentifier=plane-admin
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

[Install]
WantedBy=multi-user.target
