# deploy/packer/svc/admin.service
[Unit]
Description=Plane Admin Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/admin
EnvironmentFile=/opt/plane/admin/admin.env
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && nvm use 20 && node server.js'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
