# deploy/packer/svc/api.service
[Unit]
Description=Plane API Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/backend
EnvironmentFile=/opt/plane/backend/backend.env
ExecStart=/opt/plane/backend/bin/docker-entrypoint-api.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
