# deploy/packer/svc/api.service
[Unit]
Description=Plane API Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/backend
EnvironmentFile=/opt/plane/backend/backend.env
ExecStartPre=/bin/bash -c 'cd /opt/plane/backend && source .venv/bin/activate'
ExecStart=/bin/bash -c 'cd /opt/plane/backend && source .venv/bin/activate && exec bin/docker-entrypoint-api.sh'
Restart=always
RestartSec=10

# Logging configuration
StandardOutput=append:/var/log/plane/api.log
StandardError=append:/var/log/plane/api.error.log
SyslogIdentifier=plane-api
# Ensure log rotation
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

[Install]
WantedBy=multi-user.target
