# deploy/packer/svc/beat-worker.service
[Unit]
Description=Plane Beat Worker Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/backend
EnvironmentFile=/opt/plane/backend/backend.env
ExecStartPre=/bin/bash -c 'cd /opt/plane/backend && source .venv/bin/activate'
ExecStart=/bin/bash -c 'cd /opt/plane/backend && source .venv/bin/activate && exec bin/docker-entrypoint-beat.sh'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
