# deploy/packer/svc/space.service
[Unit]
Description=Plane Space Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/space
EnvironmentFile=/opt/plane/space/space.env
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && nvm use 20 && node space/server.js'
Restart=always
RestartSec=10

# Logging configuration
StandardOutput=append:/var/log/plane/space.log
StandardError=append:/var/log/plane/space.error.log
SyslogIdentifier=plane-space
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

[Install]
WantedBy=multi-user.target
