# deploy/packer/svc/web.service
[Unit]
Description=Plane Web Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/web
EnvironmentFile=/opt/plane/web/web.env
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && nvm use 20 && node web/server.js'
Restart=always
RestartSec=10

# Logging configuration
StandardOutput=append:/var/log/plane/web.log
StandardError=append:/var/log/plane/web.error.log
SyslogIdentifier=plane-web
LogRateLimitIntervalSec=0
LogRateLimitBurst=0

[Install]
WantedBy=multi-user.target
