# deploy/packer/svc/web.service
[Unit]
Description=Plane Web Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/plane/web
EnvironmentFile=/opt/plane/web/web.env
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && nvm use 20 && node server.js'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
